# scripts/com.ethan.playhq-updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ethan.playhq-updater</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/openclawbuddy/Documents/Codex/Ethan website/ethans-site-repo/scripts/run-weekly-playhq-update.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>3</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/openclawbuddy/Documents/Codex/Ethan website/ethans-site-repo</string>

  <key>StandardOutPath</key>
  <string>/Users/openclawbuddy/Documents/Codex/Ethan website/ethans-site-repo/logs/playhq-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/openclawbuddy/Documents/Codex/Ethan website/ethans-site-repo/logs/playhq-launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
